# src/etc/supervisord.d/ibgateway.conf
[program:ibgateway]
command=/bin/ibgateway-ctl.sh
autostart=true
autorestart=true
user=user